# wndrecd.service
[Unit]
Description=Activate Window record

[Service]
ExecStart=/usr/bin/wndrecd ~/.wndrecd.txt

[Install]
WantedBy=multi-user.target